# hugegraph-dist/src/assembly/static/bin/hugegraph.service
[Unit]
Description=HugeGraph RESTful service

[Service]
Type=forking
ExecStart=
ExecStop=

[Install]
WantedBy=multi-user.target
